# docker-compose/php/queue.ini
[program:queue]
command=php /var/www/think queue:listen --queue
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/queue_stdout.log
stderr_logfile=/var/log/queue_stderr.log