# docker-compose/php/work.ini
[program:work]
command=php /var/www/think workerman start
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/worker_stdout.log
stderr_logfile=/var/log/worker_stderr.log